# scripts/com.tradedan.deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Plantilla launchd para auto-deploy en el Mac de producción (Tradedan).

        Instalación (una sola vez en el Mac de producción):
          cp scripts/com.tradedan.deploy.plist ~/Library/LaunchAgents/
          # Editar la copia para ajustar la ruta del proyecto (ProgramArguments)
          launchctl load ~/Library/LaunchAgents/com.tradedan.deploy.plist

        Desactivar:
          launchctl unload ~/Library/LaunchAgents/com.tradedan.deploy.plist
    -->
    <key>Label</key>
    <string>com.tradedan.deploy</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Ajustar esta ruta al directorio real del proyecto en el Mac de producción -->
        <string>/Users/USUARIO/tradedan/scripts/deploy.sh</string>
    </array>

    <!-- Ejecutar cada 300 segundos (5 minutos) -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/USUARIO/tradedan/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USUARIO/tradedan/logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
